# .builder/rootfs-overlay/home/opencal/.config/systemd/user/opencal.path
[Unit]
Description=Start OpenCAL when Wayland compositor is ready

[Path]
PathExists=%t/wayland-1

[Install]
WantedBy=default.target
